# local/agency-watcher.timer
[Unit]
Description=agency watcher — every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
# Spread load a little; exact cadence doesn't matter for a health check.
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
